# src/test/testData/parser/MixedFile.cron
# Variables section
PWD = '/my/dir'

# Run every 4 minutes
*/4 * * * * mkdir "dirname" -la
@daily echo hello